# scripts/runpod/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ── vLLM ─────────────────────────────────────────────────────────────────────
; Starts first. nginx + fastify will keep retrying until vLLM is ready.
[program:vllm]
command=python3 -m vllm.entrypoints.openai.api_server
    --model %(ENV_VLLM_MODEL)s
    --host 0.0.0.0
    --port 8000
    --device %(ENV_VLLM_DEVICE)s
    --dtype %(ENV_VLLM_DTYPE)s
    --max-model-len %(ENV_VLLM_MAX_MODEL_LEN)s
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Fastify proxy ─────────────────────────────────────────────────────────────
[program:fastify]
command=/app/start-fastify.sh
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
